# dot_config/systemd/user/wttrcache.service
[Unit]
Description=Update wttr.in cache

[Service]
Type=oneshot
ExecStart=%h/.config/wttrcache.sh
